# monitoring/systemd/vpn-backup.timer
[Unit]
Description=Ежедневный бэкап vpn-selfhost (FR-8.1)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=vpn-backup.service

[Install]
WantedBy=timers.target
