# monitoring/systemd/vpn-backup.service
# Плановый бэкап (FR-8.1). Шифрует секреты; пишет в BACKUP_DIR.
[Unit]
Description=vpn-selfhost бэкап

[Service]
Type=oneshot
WorkingDirectory=/opt/vpn-selfhost
ExecStart=/usr/bin/env python3 -m vpnctl backup
