# monitoring/systemd/vpn-ru-check.service
# Проба доступности из РФ + заморозки (FR-7.1/7.2). Запускается с машины ОПЕРАТОРА
# (или любого хоста вне РФ) по таймеру; алертит в Telegram при сбое.
[Unit]
Description=vpn-selfhost проверка доступности из РФ

[Service]
Type=oneshot
WorkingDirectory=/opt/vpn-selfhost
ExecStart=/usr/bin/env python3 -m vpnctl monitor ru
ExecStartPost=/usr/bin/env python3 -m vpnctl monitor freeze
