# monitoring/systemd/vpn-healthcheck.timer
[Unit]
Description=Запускать healthcheck vpn-selfhost каждые 5 минут (FR-7)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=vpn-healthcheck.service

[Install]
WantedBy=timers.target
