# monitoring/systemd/vpn-healthcheck.service
# Локальный healthcheck на сервере (FR-7.3). Положить в /etc/systemd/system/ и
# включить парный .timer. Ненулевой выход -> журнал фиксирует сбой.
[Unit]
Description=vpn-selfhost локальный healthcheck

[Service]
Type=oneshot
User=vpnadmin
WorkingDirectory=/opt/vpn-selfhost
ExecStart=/usr/bin/env python3 -m vpnctl monitor health
